# build/bs2pro-controller.service
[Unit]
Description=BS2PRO Controller Core Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=%H/.local/bin/BS2PRO-Core --autostart
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
